# docker/etc/supervisor/conf.d/mlapi.conf
[program:mlapi]
command=python3 /mlapi/mlapi.py -c /etc/zm/mlapiconfig.ini
autostart=true
autorestart=false
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0