# docker/etc/supervisor/conf.d/mlapi.conf
[program:zoneminder]
directory=/mlapi
command=python3 ./mlapi.py -c /etc/zm/mlapiconfig.ini
autostart=true
autorestart=true
redirect_stderr=true
