# docker/etc/supervisor/conf.d/zoneminder.conf
[program:zoneminder]
command=/usr/local/bin/entrypoint.sh
autostart=true
autorestart=true
redirect_stderr=true
